# wechat_robot_test.ini
[program:wechat_robot_test]
command = python /home/kk17/Cantonese_test/wechat.py
user = kk17
redirect_stderr = true
stdout_logfile = /home/kk17/Cantonese_test/logs/robot.log

