# wechat_robot.ini
[program:wechat_robot]
command = python /home/kk17/Cantonese/wechat.py -e Pro
user = kk17
redirect_stderr = true
stdout_logfile = /home/kk17/Cantonese/logs/robot.log
